# aqt.cron
*/10 * * * * /smartmet/run/data/aqt/bin/update.sh
